# files/etc/supervisor.d/ssh.ini
[program:ssh]
command=/usr/bin/ssh -Ng -L %(ENV_PORT_EXPOSED)s:%(ENV_ADDR_TO)s:%(ENV_PORT_TO)s %(ENV_ADDR_TUNNEL)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0
